# .config/systemd/user/sunrise-sunset.timer
[Unit]
Description=update sunrise and sunset times

[Timer]
OnCalendar=hourly

[Install]
WantedBy=default.target
