# .config/systemd/user/sunrise-sunset.service
[Unit]
Description=update sunrise and sunset times

[Service]
Type=oneshot
ExecStart=/home/bryon/.local/bin/sunrise-sunset
